# rootfs/usr/lib/systemd/user/gamescope-dbus.service
[Unit]
Description=Gamescope DBus Service
After=graphical-session.target

[Service]
ExecStart=/usr/bin/gamescope-dbus
